# configs/systemd-host/cmux.target
[Unit]
Description=Cmux Host Services Target
After=network-online.target
Wants=network-online.target
Requires=cmux-openvscode.service cmux-worker.service

[Install]
WantedBy=multi-user.target
